# audit.rules
## This file is automatically generated from /etc/audit/rules.d
-D
-b 8192
-f 1
-a exit,always -F arch=b64 -S execve
-w /etc/sudoers -p wa -k fileaud_sudoers
# System files auditing
# User and group info
-w /etc/group -p a -k fileaud_sysfile_/etc/group
-w /etc/passwd -p a -k fileaud_sysfile_/etc/passwd
-w /etc/gshadow -p a -k fileaud_sysfile_/etc/gshadow
-w /etc/shadow -p a -k fileaud_sysfile_/etc/shadow
-w /etc/security/opasswd -p a -k fileaud_sysfile_/etc/security/opasswd
# Networking
-w /etc/issue -p a -k fileaud_sysfile_/etc/issue
-w /etc/issue.net -p a -k fileaud_sysfile_/etc/issue.net
-w /etc/hosts -p a -k fileaud_sysfile_/etc/hosts
-w /etc/sysconfig/network -p a -k fileaud_sysdir_/etc/sysconfig/network
-w /etc/network -p a -k fileaud_sysdir_/etc/network
# Mandatory access controls
-w /etc/selinux -p a -k fileaud_sysdir_/etc/selinux
# change uid proccess
-a always,exit -F arch=b64 -S setuid -S setgid -S setreuid -S setregid -k fileaud__changeuid
-a always,exit -F arch=b32 -S setuid -S setgid -S setreuid -S setregid -k fileaud__changeuid
# System programs and utils
-w /usr/bin -p a -k fileaud_sysdir_/usr/bin
-w /usr/local/bin -p a -k fileaud_sysdir_/usr/local/bin
-w /usr/sbin -p a -k fileaud_sysdir_/user/sbin
-w /bin -p a -k fileaud_sysdir_/bin
-w /etc/fstab -p a -k fileaud_sysfile_/etc/fstab
-w /etc/mstab -p a -k fileaud_sysfile_/etc/mstab
-w /etc/inittab -p a -k fileaud_sysfile_/etc/inittab
-w /etc/rc.local -p a -k fileaud_sysfile_/etc/rc.local
-w /etc/rc.d -p a -k sysfile_perm_edit
# Monitor config files changes (directories)
-w /etc/init.d/ -p w -k fileaud_dconf_ddir_/etc/init.d/
-w /etc/sysconfig/ -p w -k fileaud_dconf_ddir_/etc/sysconfig/
-w /etc/init/ -p w -k fileaud_dconf_ddir_/etc/init/
-w /etc/systemd/ -p w -k fileaud_dconf_ddir_/etc/systemd/
-w /boot -p wa -k fileaud_boot_chg_/boot
-w /etc/passwd -p r -F auid>=500 -F auid!=-1 -k fileaud_user_read_sysconf_/etc/passwd
-w /boot/grub/ -p r -F auid>=500 -F auid!=-1 -k fileaud_user_read_sysconf_/boot/grub/
-w /etc/ssh/ -p r -F auid>=500 -F auid!=-1 -k fileaud_user_read_sysconf_/etc/ssh/
-w /etc/init.d/ -p r -F auid>=500 -F auid!=-1 -k fileaud_user_read_sysconf_/etc/init.d/
-w /etc/fstab -p r -F auid>=500 -F auid!=-1 -k fileaud_user_read_sysconf_/etc/fstab
-w /etc/shadow -p r -F auid!=-1 -k fileaud_user_read_sysconf_/etc/shadow
# Web utilities usage audit
-w /bin/curl -p x -k fileaud_dlutil_/bin/curl
-w /bin/sftp -p x -k fileaud_dlutil_/bin/sftp
-w /usr/bin/curl -p x -k fileaud_dlutil_/usr/bin/curl
-w /usr/bin/sftp -p x -k fileaud_dlutil_/usr/bin/sftp
# access and edit data in folder
-w /home -p r -k fileaud_home_read
-w /home -p wa -k fileaud_home_mod
-w /root -p r -k fileaud_root_read
-w /root -p wa -k fileaud_root_mod
